# ansible/roles/hosts/files/etc/local.d/ansible.start
#!/bin/sh
cd /src && make ansible-pull > /var/log/ansible.log &
