# ansible/roles/hosts/files/etc/local.d/ansible.start
#!/bin/sh
cd /root && make ansible-pull &
